# filebeat.service
[Unit]
Description=Filebeat service

[Service]
User=root
WorkingDirectory=$PWD
ExecStart=$EXE
Restart=always

[Install]
WantedBy=multi-user.target
