# metricbeat.service
[Unit]
Description=Metricbeat service

[Service]
User=root
WorkingDirectory=$PWD
ExecStart=$EXE
Restart=always

[Install]
WantedBy=multi-user.target
